# go-trader.service
[Unit]
Description=Go Trading Bot (Spot + Options)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/workspace/go-trader
ExecStart=/root/.openclaw/workspace/go-trader/scheduler/trading-scheduler
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
