# go-trader.service
[Unit]
Description=Go Trading Bot (Spot + Options)
After=network.target

[Service]
Type=simple
# Run as a dedicated non-root user.
# Create before deploying: useradd --system --no-create-home go-trader
User=go-trader
Group=go-trader
WorkingDirectory=/opt/go-trader
ExecStart=/opt/go-trader/go-trader
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Prefer loading the Discord token from an env file rather than embedding
# it in the unit file or config.json:
#   echo "DISCORD_BOT_TOKEN=Bot your-token-here" > /opt/go-trader/.env
#   chmod 600 /opt/go-trader/.env
EnvironmentFile=-/opt/go-trader/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# Allow writes only to the data directory
ReadWritePaths=/opt/go-trader/scheduler

[Install]
WantedBy=multi-user.target
